# containers/apps/sogo/etc/supervisor/conf.d/supervisor-app.conf
# each program entry below is a separate terminal command.
# Each command MUST run in the foreground and stay running.
# If the command ever exits, the supervisor daemon will automatically run it again.
# Programs can be controlled like this: 'supervisorctl start fail2ban' 'supervisorctl stop fail2ban'
# supervisor writes program statuses in /var/log/supervisor

[supervisord]
nodaemon=true

[program:apache2]
startsecs=0
autostart=true
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=APACHE_ARGUMENTS="-DNO_DETACH"
command=/usr/sbin/apache2ctl start

[program:sogo]
startsecs=0
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=sogo
command=/usr/sbin/sogod -WONoDetach YES -WOPidFile /var/run/sogo/sogo.pid
