# archive/jailman/blueprints/bitwarden/includes/bitwarden.rc
#!/bin/sh

# PROVIDE: bitwarden
# REQUIRE: LOGIN DAEMON NETWORKING FILESYSTEMS
# KEYWORD: jail rust

. /etc/rc.subr

name="bitwarden"


rcvar=${name}_enable
pidfile="/var/run/${name}.pid"
command="/usr/sbin/daemon"
command_args="-u bitwarden -c -f -P ${pidfile} -r /usr/local/share/bitwarden/bin/bitwarden_rs"
load_rc_config $name
run_rc_command "$1"
